# pi/systemd/ninja-hub.service
[Unit]
Description=Ninja Command Center — multi-screen desk companion
After=network.target sound.target

[Service]
Type=simple
User=command-center
WorkingDirectory=/home/command-center/ninja-command-center/pi
ExecStartPre=/bin/bash -c 'echo 0 > /sys/class/graphics/fbcon/cursor_blink; echo -e "\033[?25l" > /dev/tty1; echo 0 > /sys/class/vtconsole/vtcon1/bind 2>/dev/null || true'
ExecStart=/usr/bin/node src/orchestrator.js
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
Environment=NODE_ENV=production
Environment=SDL_FBDEV=/dev/fb0
Environment=SDL_VIDEODRIVER=fbcon
EnvironmentFile=/home/command-center/ninja-command-center/pi/.env

[Install]
WantedBy=multi-user.target
